# deploy/marx-search-journal-send.timer
[Unit]
Description=Marx Search Journal Alerts send check (hourly; worker gates on send day + console send_time)

[Timer]
OnBootSec=8min
OnCalendar=*-*-* *:00:00 UTC
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
